# scripts/com.rivieramayapulse.upm-archive.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  LaunchAgent: archiva mensualmente el Cuadro 1.3.1 de UPM (SEGOB).
  Corre el día 15 a las 09:00 (margen sobre la publicación ~día 7 del mes M+2).
  Instalar:  cp scripts/com.rivieramayapulse.upm-archive.plist ~/Library/LaunchAgents/
             launchctl load -w ~/Library/LaunchAgents/com.rivieramayapulse.upm-archive.plist
  Quitar:    launchctl unload -w ~/Library/LaunchAgents/com.rivieramayapulse.upm-archive.plist
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.rivieramayapulse.upm-archive</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/midnightmacbookair/serie-trafico-aereo/scripts/archive_upm_1.3.1.sh</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Day</key>
    <integer>15</integer>
    <key>Hour</key>
    <integer>9</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>/Users/midnightmacbookair/serie-trafico-aereo/data/upm/launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/midnightmacbookair/serie-trafico-aereo/data/upm/launchd.err.log</string>
</dict>
</plist>
